# infra/production/systemd/fiscal-disk-check.service
[Unit]
Description=Check Fiscal production disk usage
OnFailure=fiscal-notify@%n.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/fiscal/current/infra/production/scripts/disk-check.sh
UMask=0077
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/fiscal/operations
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
LockPersonality=true
CapabilityBoundingSet=CAP_DAC_OVERRIDE
