# infra/production/systemd/fiscal-restore-verify.service
[Unit]
Description=Verify the newest Fiscal backup by isolated restore
After=postgresql.service fiscal-backup.service
Requires=postgresql.service
OnFailure=fiscal-notify@%n.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/fiscal/current/infra/production/scripts/restore-verify.sh --apply
UMask=0077
Nice=10
IOSchedulingClass=idle
# The root wrapper must drop to postgres for the isolated database lifecycle.
NoNewPrivileges=false
PrivateDevices=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/fiscal/operations
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
LockPersonality=true
CapabilityBoundingSet=CAP_SETGID CAP_SETUID CAP_DAC_OVERRIDE
