# infra/production/systemd/fiscal-backup.service
[Unit]
Description=Back up the Fiscal PostgreSQL database
After=postgresql.service
Requires=postgresql.service
OnFailure=fiscal-notify@%n.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/fiscal/current/infra/production/scripts/backup.sh --apply
UMask=0077
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/fiscal/backups /var/lib/fiscal/operations
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
LockPersonality=true
CapabilityBoundingSet=CAP_SETGID CAP_SETUID CAP_DAC_OVERRIDE
